# configs/cae.conf
[program:cae]
command=/opt/cae/start.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/cae.err.log
stdout_logfile=/var/log/cae.out.log