# configs/caefrontend.conf
[program:caefrontend]
directory=/source/CAE-Frontend/widgets/
command=grunt connect
autostart=true
autorestart=true
stderr_logfile=/var/log/cae.err.log
stdout_logfile=/var/log/cae.out.log