# configs/debugserver.conf
[program:debugserver]
command=/opt/debugserver/start.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/debugserver.err.log
stdout_logfile=/var/log/debugserver.out.log