# configs/jenkins.conf
[program:jenkins]
command=/opt/jenkins/start.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/jenkins.err.log
stdout_logfile=/var/log/jenkins.out.log